# wis2-grep-management/docker/wis2-grep-management.cron
19 * * * * su -c "wis2-grep clean --hours=$WIS2_GREP_MESSAGE_RETENTION_HOURS" > /proc/1/fd/1 2>/proc/1/fd/2
